# cron/cron.monthly/run_monthly
#!/bin/sh

LOG_DIR=/data/logs
LOG_FILE=${LOG_DIR}/cron.log
mkdir -p ${LOG_DIR}

START=$(date '+%Y-%m-%d %H:%M:%S')
echo "[${START}] [monthly] cache cleanup START" >> ${LOG_FILE}

rm -rf /data/InStock/instock/cache/hist/*
EXIT_CODE=$?

END=$(date '+%Y-%m-%d %H:%M:%S')
if [ ${EXIT_CODE} -eq 0 ]; then
    echo "[${END}] [monthly] cache cleanup END (success, started at ${START})" >> ${LOG_FILE}
else
    echo "[${END}] [monthly] cache cleanup FAILED (exit_code=${EXIT_CODE}, started at ${START})" >> ${LOG_FILE}
fi
