# cron/cron.hourly/run_hourly
#!/bin/sh

cd /data/InStock
exec uv run --frozen instock-job basic_data_daily_job
#mkdir -p /data/logs
#DATE=`date +%Y-%m-%d:%H:%M:%S`
#echo $DATE >> /data/logs/hourly.log
